# host/systemd/desktop-runner.path
[Unit]
Description=Watch the desktop-runner queue for new jobs

[Path]
PathExistsGlob=%h/desktop-runner/queue/*.job
Unit=desktop-runner.service

[Install]
WantedBy=default.target
